# deploy/systemd/dreamwedds-video-backend.service
[Unit]
Description=InvitaWedds Wedding Video API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=invitawedds
Group=invitawedds
WorkingDirectory=/var/www/invitawedds/WeddingVideoApp/backend
EnvironmentFile=/etc/invitawedds/backend.env
ExecStart=/var/www/invitawedds/WeddingVideoApp/backend/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
