# deploy/systemd/instawedds-render.service
[Unit]
Description=InvitaWedds Remotion Render Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=invitawedds
Group=invitawedds
WorkingDirectory=/var/www/invitawedds/WeddingVideoApp/render-service
EnvironmentFile=/etc/invitawedds/render.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target
